# cron.d/hp-ipmi
* * * * * root /usr/sbin/ipmi-sensors --no-header-output --comma-separated-output --entity-sensor-names > /tmp/ipmi.txt
